# packaging/debian/prerm
#!/bin/sh
set -e
systemctl stop frostd || true
systemctl disable frostd || true
